# debian/prerm
#!/bin/bash
#DEBHELPER#
set -e

#. /var/lib/dpkg/info/tm-manifest.post-args-parser $1
TMMS_SRC=/usr/lib/python3/dist-packages/tmms
TMMS_CFG_DIR=/etc/tmms

IsInstall=false
IsAbortInstall=false
IsConfigure=false

IsRemove=false
IsPurge=false
IsAbortRemove=false


for arg in "$@"
do
case $arg in
    install)
        IsInstall=true
        shift
    ;;
    abort-install)
        IsAbortInstall=true
        shift
    ;;
    configure)
        IsConfigure=true
        shift
    ;;
    remove)
        IsRemove=true
        shift
    ;;
    purge)
        IsPurge=true
        shift
    ;;
    abort-remove)
        IsAbortRemove=true
        shift
    ;;
    *)
    ;;
esac
done



if [ "$IsRemove" ]; then

    bin/systemctl stop tm-manifest-server
    bin/systemctl disable tm-manifest-server

    rm /lib/systemd/system/tm-manifest-server.service
    rm /usr/default/tm-manifest-server

    systemctl daemon-reload
fi


if [ "$IsUpgrade" ]; then
    mv $TMMS_SRC $TMMS_SRC.backup   # Backup the original tmms dir (to restore on fail)
fi
