# debian/prerm
#!/bin/bash
#DEBHELPER#

set -e	# avoids lintian error but passes it on to apt-get :-(

set -u

# Using false and true allows you to just "run the variable"
IsUpgrade=false
IsRemove=false
IsPurge=false
IsFailedUpgrade=false
IsAbortRemove=false


for arg in "$@"; do
    case $arg in
    upgrade)
        IsUpgrade=true
    ;;
    remove)
        IsRemove=true
    ;;
    purge)
        IsPurge=true
    ;;
    failed-upgrade)
        IsFailedUpgrade=true
    ;;
    abort-remove)
        IsAbortRemove=true
    ;;
    *)
    ;;
    esac
    shift
done


# Flag to indicate server is running as a daemon
if pgrep tm-manifest-server; then
    # Will flag the restart of the service later on during the update
    #touch /tmp/tm-manifest-server.restart
    echo " --- !!! WARNING !!! --- "
    echo " --- Shutting down manifesting daemon.... it will not be restarted on update! --- "
    echo " --- A systemctl service of the manifesting server will be started instead! --- "
    tm-manifest-server --daemon-stop || true
fi

# Flag to indicate server is running as systemctl service
if systemctl status tm-manifest-server | grep -q running; then
    #touch /tmp/tm-manifest-server.restart-service
    echo " --- !!! WARNING !!! --- "
    echo " --- systemctl service will be stopped and restarted on update! --- "
    systemctl stop tm-manifest-server || true
fi

if $IsRemove; then

    systemctl disable tm-manifest-server

    rm -f /lib/systemd/system/tm-manifest-server.service
    rm -f /etc/default/tm-manifest-server

    systemctl daemon-reload

    # Remove restart file, since it is not needed due to the full remove of
    # the package
    #[ -f /tmp/tm-manifest-server.restart ] && rm /tmp/tm-manifest-server.restart
    #[ -f /tmp/tm-manifest-server.restart-service ] && rm /tmp/tm-manifest-server.restart-service
fi


$IsUpgrade && echo " --- prerm upgrade: doing nothing... --- " || true

if $IsFailedUpgrade; then
    echo " --- prerm failed-upgrade: doing nothing... (handled by postrm) --- "
fi


if $IsAbortRemove; then
    echo " --- prerm abort-remove: is not handlede... skipping. ---"
fi

exit 0
